# .config/systemd/user/polybar.service
[Unit]
Description=My desktop bar

[Service]
Type=simple
ExecStart=/home/icelk/scripts/polybar.sh
