# config/systemd/backup.service
[Unit]
Description=Backup and snapshot

[Service]
Type=simple
ExecStart=sh -c 'sleep 30; notify-send -a "Backup" "Starting..."; btrbk run; notify-send -a "Backup" Complete!'
