# config/systemd/backup.service
[Unit]
Description=Backup and snapshot

[Service]
Type=simple
ExecStart=btrbk run
