# config/systemd/check-failed.service
[Unit]
Description=Check systemd failed units

[Service]
Type=oneshot
ExecStart=/home/icelk/scripts/check-failed.sh

