# config/systemd/clipcat.service
[Unit]
Description=Starts the clipboard daemon

[Service]
Type=simple
ExecStart=clipcatd --no-daemon

[Install]
WantedBy=default.target
