# config/systemd/dunst.service
[Unit]
Description=Notifications

[Service]
Environment=WAYLAND_DISPLAY=wayland-0
Type=simple
ExecStart=dunst
