# config/systemd/dunst.service
[Unit]
Description=Notifications

[Service]
Type=simple
ExecStart=dunst
