# config/systemd/nitrogen.service
[Unit]
Description=Start (restore) wallpaper

[Service]
Type=oneshot
ExecStart=nitrogen --restore
Restart=on-failure
RestartSec=5
