# config/systemd/notify-reminder.service
[Unit]
Description=Send notifications from each line in ~/Reminder.txt

[Service]
Type=simple
ExecStart=sh -c "~/scripts/send-reminders.sh"
