# config/systemd/notify-reminder.timer
[Unit]
Description=Reminders to check

[Timer]
OnCalendar=Sat *-*-* 14:00:00
Persistent=true

[Install]
WantedBy=timers.target
