# config/systemd/notify-reminder.service
[Unit]
Description=Send notification with a document

[Service]
Type=oneshot
ExecStart=/home/icelk/scripts/send-reminders.sh

