# config/systemd/picom.service
[Unit]
Description=Composter

[Service]
Type=simple
ExecStart=picom --experimental-backends

