# config/systemd/polybar.service
[Unit]
Description=My desktop bar

[Service]
Type=simple
ExecStart=/home/icelk/scripts/polybar.sh
ExecStop=rm /tmp/run/packages-upgradable-count.pid
