# config/systemd/polybar.service
[Unit]
Description=My desktop bar

[Service]
Type=simple
ExecStart=sh -c "~/scripts/polybar.sh"
ExecStop=rm /tmp/run/packages-upgradable-count.pid
