# config/systemd/sync.service
[Unit]
Description=Sync
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=sh -c "~/scripts/rsync.sh"
