# config/systemd/sync.timer
[Unit]
Description=Sync to nas every day

[Timer]
OnCalendar=17:00
Persistent=true

[Install]
WantedBy=timers.target

