# config/systemd/xsnow.service
[Unit]
Description=Set the weather to snow.

[Service]
ExecStart=sh -c "~/scripts/xsnow.sh"
